# chef-backend.service
[Unit]
Description=Chef Anti-Inflatie Backend API
After=network.target

[Service]
Type=simple
User=arch
WorkingDirectory=/home/arch/IOferta/backend
Environment="PATH=/home/arch/IOferta/backend/venv/bin"
ExecStart=/home/arch/IOferta/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
